# mysql_threads.crontab
*/10 * * * * /usr/local/sbin/mysql_threads.sh
